# scripts/runtime/com.arvectum.embeddings.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict><key>Label</key><string>com.arvectum.embeddings</string><key>ProgramArguments</key><array><string>/opt/homebrew/bin/llama-server</string><string>--model</string><string>/Users/master/models/embeddings/Qwen3-Embedding-4B-GGUF/Qwen3-Embedding-4B-Q8_0.gguf</string><string>--alias</string><string>Qwen3-Embedding-4B</string><string>--embedding</string><string>--pooling</string><string>last</string><string>--ctx-size</string><string>8192</string><string>--parallel</string><string>1</string><string>--batch-size</string><string>2048</string><string>--ubatch-size</string><string>2048</string><string>--threads</string><string>8</string><string>--threads-batch</string><string>8</string><string>--no-cont-batching</string><string>--host</string><string>127.0.0.1</string><string>--port</string><string>8090</string></array><key>RunAtLoad</key><true/><key>KeepAlive</key><true/><key>StandardOutPath</key><string>__LOG_DIR__/embeddings.log</string><key>StandardErrorPath</key><string>__LOG_DIR__/embeddings.error.log</string></dict></plist>
