# scripts/runtime/com.arvectum.backend.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict><key>Label</key><string>com.arvectum.backend</string><key>ProgramArguments</key><array><string>/bin/bash</string><string>__ROOT__/scripts/runtime/start.sh</string></array><key>WorkingDirectory</key><string>__ROOT__</string><key>RunAtLoad</key><true/><key>KeepAlive</key><false/><key>StandardOutPath</key><string>__LOG_DIR__/backend.log</string><key>StandardErrorPath</key><string>__LOG_DIR__/backend.error.log</string></dict></plist>
